# etc/rc.local
/usr/local/sbin/start_hostapd.sh
/usr/local/sbin/set_on_boot.sh
